# updatedns.service
[Unit]
Description=Update ClouDNS
After=network.target

[Service]
ExecStart=/usr/local/sbin/update_cloudns.py

[Install]
WantedBy=multi-user.target
